# com.tradebot.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Label -->
    <key>Label</key>
    <string>com.tradebot.agent</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/rdovgan/Desktop/projects/trade-bot/scripts/start_bot.sh</string>
    </array>

    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/rdovgan/Desktop/projects/trade-bot</string>

    <!-- Run at load (when system starts) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Keep alive - restart if it crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Start interval - run daily at 8:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/rdovgan/Desktop/projects/trade-bot/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rdovgan/Desktop/projects/trade-bot/logs/launchd_stderr.log</string>

    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Process Type -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Throttle restart attempts -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
