# tradebot.timer
[Unit]
Description=Trading Bot Daily Timer
Requires=tradebot.service

[Timer]
# Run daily at 8:00 AM UTC
OnCalendar=daily
OnCalendar=*-*-* 08:00:00

# Run 5 minutes after boot if missed
Persistent=true

[Install]
WantedBy=timers.target
