# tradebot.service
[Unit]
Description=Trading Bot Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rdovgan
Group=rdovgan
WorkingDirectory=/home/rdovgan/trade-bot
Environment="PATH=/home/rdovgan/trade-bot/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# Pre-start: Update code from git
ExecStartPre=/bin/bash -c 'cd /home/rdovgan/trade-bot && git pull origin master'

# Start the bot
ExecStart=/home/rdovgan/trade-bot/.venv/bin/python -m trade_bot.main

# Restart behavior
Restart=on-failure
RestartSec=60
StartLimitBurst=5
StartLimitIntervalSec=300

# Logging
StandardOutput=append:/home/rdovgan/trade-bot/logs/bot.log
StandardError=append:/home/rdovgan/trade-bot/logs/bot_error.log

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/rdovgan/trade-bot/logs /home/rdovgan/trade-bot/data

[Install]
WantedBy=multi-user.target
