# tradebot.service
[Unit]
Description=Trading Bot
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/trade-bot
Environment="PATH=/root/trade-bot/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# Git pull before starting
ExecStartPre=/bin/bash -c 'cd /root/trade-bot && git pull origin master'

# Run the bot
ExecStart=/root/trade-bot/.venv/bin/python -m trade_bot.main

# Auto-restart on failure
Restart=on-failure
RestartSec=60

# Logging
StandardOutput=append:/root/trade-bot/logs/bot.log
StandardError=append:/root/trade-bot/logs/bot.log

[Install]
WantedBy=multi-user.target
